# crontab
0 12 * * * /Users/shinkeika/anaconda3/bin/python3.7 /Users/shinkeika/pythonproject/setwallpaper/setwallpapereveryday.py 